# tools/target/linux/ramips/base-files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
#ttyS0::askfirst:/bin/ash --login
#ttyS1::askfirst:/bin/ash --login
::respawn:/sbin/getty -L ttyS1 115200 vt100
